# deploy/dailyhot-api-on-failure.service
[Unit]
Description=DailyHotApi failure logger

[Service]
Type=oneshot
ExecStart=/usr/local/lib/dailyhot/on-failure.sh
